# profile.d/hadoop-path.sh
export PATH=/opt/hadoop-3.2.1/bin:/opt/hadoop-3.2.1/sbin:$PATH
