# profile.d/spark-path.sh
export PATH=/spark/bin:/spark/sbin:$PATH
